# 20211022_WebIntro/3_das_blog/other/mysql.conf
[program:mysql]
command=mysqld_safe

autorestart=true
exitcodes=0

stopsignal=TERM
stopasgroup=true
stopwaitsecs=10

stdout_events_enabled=true
stderr_events_enabled=true
